# auto_update.crontab
#MIN HOUR DAY MONTH DAYWK   CMD
*/5  *    *   *     *       /home/jfrasca/repos/meshuggahme/auto_updater.sh
