# auto_update.crontab
#MIN HOUR DAY MONTH DAYWK   CMD
*/5  *    *   *     *       /home/meshuggah/meshuggahme/auto_updater.sh
